# bts_tools/templates/deploy/supervisord.conf
{% for client in config_yaml['clients'] %}

[program:{{ client }}witness]
user={{ unix_user }}
command=/home/{{ unix_user }}/.virtualenvs/bts_tools/bin/bts run {{ client }}
stopsignal=INT
stopwaitsecs=60

[program:{{ client }}cli]
user={{ unix_user }}
command=sh -c "/home/{{ unix_user }}/.virtualenvs/bts_tools/bin/bts run_cli {{ client }}; sleep 10"
autorestart=true
startsecs=0
redirect_stderr=true

{% endfor %}
