# bts_tools/templates/deploy/supervisord.conf
{% for client in config_yaml['clients'] %}

[program:{{ client }}witness]
user={{ unix_user }}
command=/home/{{ unix_user }}/.virtualenvs/bts_tools/bin/bts run {{ client }}

[program:{{ client }}cli]
user={{ unix_user }}
command=sh -c "sleep 10; /home/{{ unix_user }}/.virtualenvs/bts_tools/bin/bts run_cli {{ client }}"
autorestart=true
startsecs=0
redirect_stderr=true

{% endfor %}
