# bts_tools/templates/deploy/supervisord.conf
[program:witness]
user={{ unix_user }}
command=/home/{{ unix_user }}/.virtualenvs/bts_tools/bin/bts2 run

[program:cli]
user={{ unix_user }}
command=sh -c "sleep 10; /home/{{ unix_user }}/.virtualenvs/bts_tools/bin/bts2 run_cli"
autorestart=true
startsecs=0
redirect_stderr=true

